# systemd/reddit-digest.timer
[Unit]
Description=Run Reddit Digest every 30 minutes

[Timer]
# Run 30 minutes after boot, then every 30 minutes
OnBootSec=5min
OnUnitActiveSec=30min

# Add randomized delay to avoid thundering herd (0-60 seconds)
RandomizedDelaySec=60

# Persist timer across reboots
Persistent=true

[Install]
WantedBy=timers.target
